# Deploy/places-host/podman-compose.service
[Unit]
Description=places-host podman-compose stack (pg_places + link_extractor)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/ec2-user/ScrapingShoppingMalls/Deploy/places-host
ExecStart=/usr/local/bin/podman-compose up -d
ExecStop=/usr/local/bin/podman-compose down
TimeoutStartSec=600

[Install]
WantedBy=default.target
