# Deploy/podman-prune.timer
[Unit]
Description=Weekly podman cleanup (Sundays 03:00 UTC)

[Timer]
OnCalendar=Sun *-*-* 03:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
