# Deploy/queue-host/podman-compose.service
[Unit]
Description=queue-host podman-compose stack (pg_queue + scraper)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/SearchingThings/Deploy/queue-host
# Self-update to latest main before starting. The leading `-` makes these
# best-effort: if GitHub is unreachable at boot, the stack still starts on the
# code already on disk instead of failing. Gitignored files (.env) are untouched.
ExecStartPre=-/usr/bin/git -C %h/SearchingThings fetch --prune origin main
ExecStartPre=-/usr/bin/git -C %h/SearchingThings reset --hard origin/main
ExecStart=/usr/local/bin/podman-compose up -d --build
ExecStop=/usr/local/bin/podman-compose down
TimeoutStartSec=600

[Install]
WantedBy=default.target
